# dynamic-load-balancer/haproxy/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
numprocs = 1
startretries = 0
command=service haproxy start